# virt-who.service
[Unit]
Description=Daemon for reporting virtual guest IDs to subscription-manager
After=libvirtd.service network.target

[Service]
Type=notify
PIDFile=/run/virt-who.pid
ExecStart=/usr/bin/virt-who
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
